# nfl-bot.service
[Unit]
Description=NFL Discord Bot
After=network.target

[Service]
Type=simple
User=opc
WorkingDirectory=/opt/nfl-bot
EnvironmentFile=/etc/nfl-bot.env
ExecStart=/opt/nfl-bot/venv/bin/python bot.py
Restart=always
RestartSec=10
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
